# ccsrouter/startup-services.conf
[program:nginx]
command=/usr/sbin/nginx

[program:haproxy]
command=bash /tmp/run-haproxy.sh "%(ENV_CCSAPI_URL)s"
